# Software_RK3588/services/gimbal-ai.service
[Unit]
Description=Gimbal AI Assistant (RK3588 Headless)
After=network.target

[Service]
Type=simple
# install_service.sh will replace the placeholders below.
User={{RUN_USER}}
WorkingDirectory={{SOFTWARE_DIR}}
Environment=GIMBAL_AI_PROJECT_ROOT={{PROJECT_ROOT}}
ExecStart={{PYTHON_BIN}} {{SOFTWARE_DIR}}/main.py --headless --config {{SOFTWARE_DIR}}/configs/default.yaml
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
